# templates/default/titan.upstart.conf.erb
description "Titan Server"

author "brian.cajes@gmail.com"

version "0.4.x"

#eg. /var/log/upstart/titan.log
console log

env USER=<%= @node.titan.user %>

chdir <%= @node.titan.installation_dir %>

start on runlevel [2345]

stop on runlevel [016]

respawn

pre-start script
	  exec sudo -u <%= node.titan.user%> <%= node.titan.start_command %> 
	  sleep 5
end script
post-stop script
	  exec sudo -u <%= node.titan.user %> <%= node.titan.stop_command %> 
	  sleep 5
end script




